# raspberry/iot.service
[Unit]
Description=IOT Service

[Service]
User=pi
ExecStart=/usr/bin/python3 /home/caner/Project/iot/raspberry/main.py
Restart=always
RestartSec=60
